# test/data/deploy/production/tmp.conf
description "tmp node process"
start on runlevel [2]
stop on runlevel [016]
console owner
chdir "/var/www/tmp"
env NODE_ENV=production
exec node ./bin/paratrooper >> /var/log/tmp.log 2>&1
respawn
respawn limit 5 15
